# systemd/osint-backend.service
[Unit]
Description=OSINT Backend Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/osint-framework
ExecStart=/usr/bin/node /opt/osint-framework/server.js
Restart=always
User=www-data
Environment=NODE_ENV=production
Environment=PORT=4000

[Install]
WantedBy=multi-user.target
